# host-a/docker/keepalived-supervisord.conf
[program:keepalived]
command=/usr/sbin/keepalived
